# cron-tasks/cron.hourly/update-external-ip
#!/bin/sh

proxy=$(cat /tmp/proxy)
if [ -z "$proxy" ]; then
  curl --proxy "$proxy" -sL icanhazip.com -o /tmp/external-ip
else
  curl -sL icanhazip.com -o /tmp/external-ip
fi
