# cron-tasks/cron.hourly/update-external-ip
#!/bin/sh

curl -sL icanhazip.com -o /tmp/external-ip
