# cron-tasks/cron.hourly/update-external-ip
#!/bin/sh

if [ -z "$http_proxy" ]; then
  curl -sL icanhazip.com -o /tmp/external-ip
else
  curl --proxy "$http_proxy" -sL icanhazip.com -o /tmp/external-ip
fi
