# cron-tasks/cron.hourly/update-external-ip.sh
#!/usr/bin/env bash

curl -sL icanhazip.com -o /tmp/external-ip
